# package/luci-app-setup/files/etc/uci-defaults/40_luci-app-setup
touch /etc/config/setup
touch /etc/config/wireless

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@setup[-1]
	add ucitrack setup
	set ucitrack.@setup[-1].init=setup
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
