# package/luci-app-setup/files/etc/uci-defaults/ali_c
#!/bin/sh

command -v alemprator_c >/dev/null 2>&1 && alemprator_c &
